# indexator.service
[Unit]
Description=Indexator - FastAPI Backend
After=network.target

[Service]
Type=simple
User=indexator
Group=indexator
WorkingDirectory=/var/www/indexator
Environment="PATH=/var/www/indexator/venv/bin:/usr/bin:/bin"
Environment="PLAYWRIGHT_BROWSERS_PATH=/var/www/indexator/.cache/ms-playwright"
ExecStart=/var/www/indexator/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

# Security
NoNewPrivileges=true
PrivateTmp=true

# Limits
LimitNOFILE=65536
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
